# resources/jenkins_backup_template.timer
[Unit]
Description=Runs Jenkins backup once per day

[Timer]
OnCalendar=daily
RandomizedDelaySec=10

[Install]
WantedBy=timers.target
